# deploy/telegram-bot.service
[Unit]
Description=Interview Router - Telegram Bot
After=network.target

[Service]
Type=simple
User=elmanamador
WorkingDirectory=/home/elmanamador/interview-router
ExecStart=/home/elmanamador/interview-router/venv/bin/python telegram_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
